# units/updater.service
[Unit]
Description=CLIP OS system updater
ConditionPathExists=!/run/update_ready

[Service]
Type=simple
ExecStart=/usr/bin/updater
